# deploy/vps-dry-run/systemd/smartwork-production-worker.dry-run.service
[Unit]
Description=SmartWork Agent Production Worker Dry-Run
After=network.target smartwork-control-server.dry-run.service

[Service]
Type=simple
WorkingDirectory=/opt/smartwork-agent
Environment=NODE_ENV=production
Environment=SMARTWORK_DRY_RUN=true
Environment=SMARTWORK_NO_SIAGA_INPUT=true
Environment=SMARTWORK_NO_BROWSER_OPEN=true
Environment=SMARTWORK_NO_REAL_SAVE=true
Environment=SMARTWORK_NO_REAL_SEND=true
Environment=SMARTWORK_REAL_SAVE_ENABLED=false
Environment=SMARTWORK_WORKER_INTERVAL_MS=1000
Environment=SMARTWORK_DELIVERY_MODE=app_download_only
ExecStart=/usr/bin/node scripts/smartwork-production-worker.mjs --daemon --dry-run
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

