# deploy/dewavps/smartwork-production-worker.service

[Unit]
Description=SmartWork Production Worker Daemon
After=network.target smartwork-control-server.service
Requires=smartwork-control-server.service

[Service]
Type=simple
WorkingDirectory=/opt/smartwork-agent
EnvironmentFile=/opt/smartwork-agent/.env.production
ExecStart=/usr/bin/node scripts/smartwork-production-worker-daemon.mjs
Restart=always
RestartSec=5
User=smartwork
Group=smartwork

[Install]
WantedBy=multi-user.target
